# .pc/var-run.diff/instfiles/xrdp.service.in
[Unit]
Description=xrdp daemon
Documentation=man:xrdp(8) man:xrdp.ini(5)
Requires=xrdp-sesman.service
After=network.target xrdp-sesman.service

[Service]
Type=forking
PIDFile=@localstatedir@/run/xrdp/xrdp.pid
RuntimeDirectory=xrdp
EnvironmentFile=-@sysconfdir@/sysconfig/xrdp
EnvironmentFile=-@sysconfdir@/default/xrdp
User=xrdp
Group=xrdp
PermissionsStartOnly=true
ExecStartPre=/bin/sh /usr/share/xrdp/socksetup
ExecStart=@sbindir@/xrdp $XRDP_OPTIONS
ExecStop=@sbindir@/xrdp $XRDP_OPTIONS --kill

[Install]
WantedBy=multi-user.target
